# dot_config/systemd/user/hypr-monitors.path
[Unit]
Description=Watch for monitor configuration changes
After=graphical-session.target

[Path]
PathModified=/sys/class/drm
PathModified=/proc/modules

[Install]
WantedBy=graphical-session.target