# dot_config/systemd/user/hyprland-configreload-listener.service
[Unit]
Description=Hyprland event listener (config reload, monitor changes)
After=graphical-session.target waybar.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/scripts/hyprland-configreload-listener.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
