# dot_config/systemd/user/waybar-reload.service
[Unit]
Description=Reload Waybar on config/style changes

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill -SIGUSR2 waybar
