# dot_config/systemd/user/wayland-env-guard.service
[Unit]
Description=Repoint systemd user env at a live Wayland socket

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wayland-env-guard
