# dot_config/systemd/user/wayland-wm@.service.d/10-resource-tree.conf
[Service]
# An OOM kill of any child must not take down the compositor (this unit
# previously failed with result 'oom-kill' whenever the kernel shot a child,
# dumping the session to the greeter)
OOMPolicy=continue
# Reserve memory for the compositor tree under global pressure
MemoryLow=2G
# Compositor wins CPU contention against default-weight (100) app scopes
CPUWeight=500
