# dot_config/systemd/user/waybar.service.d/restart-on-failure.conf
[Service]
Restart=on-failure
RestartSec=2
StartLimitInterval=30
StartLimitBurst=5
ExecStart=
ExecStart=/usr/bin/waybar --log-level debug
