# dot_config/systemd/user/xdg-desktop-portal-hyprland.service.d/override.conf
[Unit]
# D-Bus activation can race a compositor restart; stock burst limit (3) makes
# xdph die permanently into start-limit-hit (2026-07: 6 days of dead screen
# sharing). Retry forever — ConditionEnvironment + PartOf=graphical-session
# keep it from spinning when there's no session at all.
StartLimitIntervalSec=0

[Service]
RestartSec=5
